# templates/gitea.service.j2
[Unit]
Description=Gitead Project Management program
After=network.target

[Service]
Type=simple
Restart=no
User=_gitea
RootDirectory=/
WorkingDirectory={{ gitea_install_path }}
ExecStart={{ gitea_install_path }}/{{ gitea_binary }} web -p {{ http_port }}


[Install]
WantedBy=multi-user.target
