# templates/gitea.service.j2
[Unit]
Description=Gitead Project Management program
After=network.target

[Service]
Type=simple
Restart=on-failure
User=_gitea
RootDirectory=/
WorkingDirectory={{ gitea_install_path }}
ExecStart={{ gitea_install_path }}/{{ gitea_wrapper }} web -p {{ gitea_http_port }} 

[Install]
WantedBy=multi-user.target
