# templates/gitea.service.j2
[Unit]
Description=Gitea Project Management program
After=network.target

[Service]
Type=simple
Restart=on-failure
User=_gitea
UMask=0077
WorkingDirectory={{ gitea_install_path }}
ExecStart={{ gitea_install_path }}/{{ gitea_binary }} web -p {{ gitea_http_port }} 

[Install]
WantedBy=multi-user.target
